# .local/share/systemd/user/KeePassXC.service
# (User) Service Unit for KeePassXC
[Unit]
Description=KeePassXC with automatic restart
Documentation=man:keepassxc-cli(1) man:keepassxc(1)
After=graphical.target

[Service]
Type=exec
WorkingDirectory=%h
ExecStart=/usr/bin/env keepassxc -platform xcb
Restart=always
RestartSec=5s

[Install]
WantedBy=default.target

